# templates/mailhog.service.erb
[Unit]
Description=MailHog Service
After=syslog.target network.service

[Service]
EnvironmentFile=<%= scope['mailhog::config'] %>
Type=simple
ExecStart=<%= scope['mailhog::binary_file'] %>
StandardOutput=journal

[Install]
WantedBy=multi-user.target
